# thingdb.service
[Unit]
Description=ThingDB Inventory Management System
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=thingdb
Group=thingdb
WorkingDirectory=/var/lib/thingdb/app
Environment="PATH=/var/lib/thingdb/app/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/var/lib/thingdb/app/venv/bin/gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 2 \
    --timeout 600 \
    --access-logfile - \
    --error-logfile - \
    "thingdb.main:app"

# Restart policy
Restart=always
RestartSec=10

# Security settings
# NoNewPrivileges disabled to allow sudo for power management
# thingdb user has NOPASSWD access only to: shutdown, reboot, sync, restart
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thingdb

[Install]
WantedBy=multi-user.target

