# pi-setup/ble-config.service
[Unit]
Description=Bluetooth LE Configuration Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/var/lib/inventory/app/venv/bin/python /usr/local/bin/ble_config_service.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
