# updates/inventory-recovery.service
[Unit]
Description=Inventory System Recovery Service
Documentation=Manages upgrade flags and automatic rollback
After=network.target
Before=inventory-app.service
Wants=inventory-app.service

[Service]
Type=simple
ExecStart=/usr/local/bin/inventory-recovery
Restart=always
RestartSec=10
User=root
Group=root

# Ensure this service starts even if others fail
RemainAfterExit=no
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
